# localProjects/cmangos_projects/transfer/cmangos-daily-sync.service
[Unit]
Description=CMaNGOS Daily Character Sync (Classic→TBC→WotLK, multi-account)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/cmangos-transfer/daily-sync.sh
EnvironmentFile=/opt/cmangos-transfer/.env
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
