# localProjects/cmangos_projects/transfer/cmangos-daily-sync.timer
[Unit]
Description=CMaNGOS Daily Character Sync Timer

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
